# systemd/camera-uploader.service
[Unit]
Description=Camera Uploader FastAPI Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=camera-uploader
Group=camera-uploader
WorkingDirectory=/opt/camera-uploader
Environment=CAMERA_UPLOADER_CONFIG=/opt/camera-uploader/config.json
Environment=PYTHONPATH=/opt/camera-uploader
ExecStart=/opt/camera-uploader/venv/bin/python -m uvicorn src.main:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/camera-uploader

[Install]
WantedBy=multi-user.target
